# lovetap.service
[Unit]
Description=LoveTap Landing Page
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/lovetap
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=PORT=3000
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

# 确保进程有足够的文件描述符限制
LimitNOFILE=65535

# 安全相关设置
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
