# lovetap.service
[Unit]
Description=LoveTap Landing Page
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/root/desktop/website/LoveTap-LandingPage
ExecStart=/usr/bin/node server.js
Restart=on-failure
Environment=PORT=3000
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
